# profile.d/deb-packaging.sh
#!/bin/sh
export DEBFULLNAME='John P. Lettman'
export DEBEMAIL='the@johnlettman.com'
export DEBSIGN_KEYID='198FCD6C917B5FFEB34B7F2F17EE3A0E7BE7193C'
export QUILT_PATCHES='debian/patches'
