# profile.d/gpg.sh
#!/bin/sh
# Use the following key for GPG: John Lettman <the@johnlettman.com>
export GPGKEY='C0CECF66705141610001CA728CE505CF86B53CC3'
