# profile.d/kwin.sh
#!/bin/sh
export KWIN_TRIPLE_BUFFER=1
export KWIN_COMPOSE=O2
