# profile.d/less.sh
#!/bin/sh
export PAGER=less
export LESS='-F -Q -M -R -X -i -g -s -x4 -z-2'
export LESSCOLORIZER='bat'
