# profile.d/snapcraft.sh
#!/bin/sh
export SNAPCRAFT_BUILD_ENVIRONMENT=lxd