# profile.d/xdg.sh
#!/bin/sh
# Fix missing snap application icons in KDE
export XDG_DATA_DIRS="${XDG_DATA_DIRS}:/var/lib/snapd/desktop"
