# shell/profile.d/lang-go
#!/bin/sh
if command -v go >/dev/null 2>&1; then
    export GOPATH="${HOME}/repos/go";
    append_path "${GOPATH}/bin";

    if [ ! -d "${GOPATH}" ]; then
        mkdir -p "${GOPATH}";
    fi;
fi;