# shell/profile.d/tool-batcat
#!/bin/sh
if program_exists batcat; then
    # Fix Ubuntu `bat` utility naming
    alias bat=batcat
fi
