# shell/profile.d/tool-gists
#!/bin/sh
export GISTS="${HOME}/repos/gists"

if [ ! -d "${GISTS}" ]; then
    mkdir -p "${GISTS}"
fi

# shellcheck disable=SC2139
alias gists="cd ${GISTS}"
