# shell/profile.d/tweak-keybinds
#!/bin/sh
# Set emacs keybindings on shells that support them.
# shellcheck disable=SC2009
if (ps -o comm= $$ | grep -Eq '^(z|ba)sh'); then
    # shellcheck disable=SC3040
    set -o emacs
fi
