# shell/profile.d/ui-lesspipe
#!/bin/sh
lesspipe="${HOME}/.config/less/lesspipe/lesspipe.sh";
if [ -f "${lesspipe}" ]; then
    eval "$(SHELL=/bin/sh "${lesspipe}")";
fi;
